# bin/numatuned.service
[Unit]
Description=Replacement of numad for Linux hypervisors running libvirt, qemu

[Service]
Type=simple
ExecStart=/usr/bin/numatuned
Restart=always

[Install]
WantedBy=multi-user.target
